# backend/securesight.conf
[program:securesight-backend]
command=/var/www/securesight-backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/var/www/securesight-backend
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/securesight-backend.log
stderr_logfile=/var/log/securesight-backend-error.log
environment=PATH="/var/www/securesight-backend/venv/bin",HOME="/var/www"

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
